# deploy/socialforge.service
[Unit]
Description=SocialForge gRPC Server
After=network.target postgresql.service

[Service]
Type=simple
User=socialforge
Group=socialforge
WorkingDirectory=/opt/socialforge
ExecStart=/opt/socialforge/bin/server
Restart=always
RestartSec=5
Environment=AUTH0_DOMAIN=your-auth0-domain
Environment=AUTH0_CLIENT_ID=your-auth0-client-id
Environment=AUTH0_CLIENT_SECRET=your-auth0-client-secret

[Install]
WantedBy=multi-user.target 